# lib/modules/manager/quadlet/__fixtures__/registry-alias.container
[Container]
Image=quay.io/metallb/controller:v0.13.10
